# 4.5.0.21/RPMs/list/busybox-symlinks-dhcp.list
/usr/lib/systemd/system/udhcpd.service
/usr/sbin/udhcpc
/usr/sbin/udhcpd
