# 4.5.0.21/RPMs/list/dbus.list
/etc/dbus-1
/etc/dbus-1/session.conf
/etc/dbus-1/session.d
/etc/dbus-1/system.conf
/etc/dbus-1/system.d
/usr/bin/dbus-cleanup-sockets
/usr/bin/dbus-daemon
/usr/bin/dbus-launch
/usr/bin/dbus-monitor
/usr/bin/dbus-run-session
/usr/bin/dbus-send
/usr/bin/dbus-test-tool
/usr/bin/dbus-update-activation-environment
/usr/bin/dbus-uuidgen
/usr/lib/systemd/system/basic.target.wants/dbus.service
/usr/lib/systemd/system/dbus.service
/usr/lib/systemd/system/dbus.socket
/usr/lib/systemd/system/dbus.target.wants/dbus.socket
/usr/lib/systemd/system/multi-user.target.wants/dbus.service
/usr/lib/systemd/system/sockets.target.wants/dbus.socket
/usr/lib/systemd/user/dbus.service
/usr/lib/systemd/user/dbus.socket
/usr/libexec/dbus-1
/usr/libexec/dbus-1/dbus-daemon-launch-helper
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/services
/usr/share/dbus-1/session.conf
/usr/share/dbus-1/session.d
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system.conf
/usr/share/dbus-1/system.d
/usr/share/licenses/dbus-1.14.4+git1
/usr/share/licenses/dbus-1.14.4+git1/COPYING
/var/lib/dbus
/var/run/dbus
