# 4.5.0.21/RPMs/list/nemo-password-manager.list
/etc/dbus-1/system.d/org.nemo.passwordmanager.conf
/usr/bin/nemo-password-manager
/usr/lib/systemd/system/dbus-org.nemo.passwordmanager.service
/usr/share/dbus-1/system-services/org.nemo.passwordmanager.service
/usr/share/licenses/nemo-password-manager-0.1.5
/usr/share/licenses/nemo-password-manager-0.1.5/COPYING
