# 4.5.0.21/RPMs/list/systemd-config-sailfish.list
/etc/systemd/bootchart.conf
/etc/systemd/journald.conf
/etc/systemd/logind.conf
/etc/systemd/system.conf
/etc/systemd/user.conf
/etc/udev/udev.conf
/usr/lib/systemd/system/default.target
