# 4.5.0.21/RPMs.data/list/droid-config-xqau52.list
/etc/acdbdata
/etc/gbinder.conf
/etc/gst-droid/gstdroidcodec.conf
/etc/libncicore.conf
/etc/mce/60-fingerprint-sensor.conf
/etc/mce/60-hybris-features.ini
/etc/mce/60-memnotify-seine.conf
/etc/mce/60-proximity-sensor.conf
/etc/modprobe.d/blacklist-evbug.conf
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/syncevolution
/etc/obexd/opp.conf
/etc/obexd/pbap.conf
/etc/obexd/root
/etc/ofono/binder.conf
/etc/ofono/binder.d/dual-sim.conf
/etc/ofono/ril_subscription.conf
/etc/ofono/ril_subscription.d/dual-sim.conf
/etc/sailfish-fpd/50-settings.ini
/etc/selinux/config
/etc/selinux/minimum/contexts/dbus_contexts
/etc/selinux/minimum/contexts/files/file_contexts
/etc/selinux/minimum/policy/policy.30
/etc/sensorfw/primaryuse.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/usb-moded/usb-moded-configfs-seine.ini
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/zypp/systemCheck.d/ha.check
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/oem
/system/lib/libandroidicu.so
/system/lib/libicui18n.so
/system/lib/libicuuc.so
/system/lib64/libandroidicu.so
/usr/bin/droid/droid-bootctl.sh
/usr/bin/droid/droid-hal-early-init.sh
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/droid-load-firmware.sh
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/wait_for_hwcomposer
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/autologin@.service.d/50-vendor-hwcomposer-2-3.conf
/usr/lib/systemd/system/droid-bootctl.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/droid-reclaim-memory.service
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/systemd/system/graphical.target.wants/droid-bootctl.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/jolla-actdead-charging.service.d/10-wait_for_hwcomposer.conf
/usr/lib/systemd/system/multi-user.target.wants/wait_for_keymaster.service
/usr/lib/systemd/system/network.target.wants/wifisetup.service
/usr/lib/systemd/system/sailfish-unlock-agent.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/system/sailfish-upgrade-ui.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/lib/systemd/system/wait_for_keymaster.service
/usr/lib/systemd/system/wifisetup.service
/usr/lib/systemd/system/yamuisplash.service.d/50-vendor-hwcomposer-2-3.conf
/usr/lib/systemd/user/jolla-startupwizard-pre-user-session.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/user/lipstick.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/libexec/droid-hybris/system/etc/init/android.hardware.media.omx@1.0-service.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/usr/libexec/droid-hybris/system/etc/init/ecclist.rc
/usr/libexec/droid-hybris/system/etc/init/hvdcp_opti.rc
/usr/libexec/droid-hybris/system/etc/ld.config.29.txt
/usr/share/csd/settings.d/10-hw-settings.ini
/usr/share/csd/settings.d/20-hw-settings-dual-sim.ini
/usr/share/ngfd/plugins.d/60-profile.ini
/usr/share/qt5/keymaps/droid.qmap
/usr/share/ssu/board-mappings.d/05-sony-seine.ini
/usr/share/ssu/board-mappings.d/10-xqau52.ini
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/nemo/99-qtscenegraph.conf
/var/lib/environment/ofono/noplugin.conf
