# 4.5.0.21/RPMs.data/list/dsme.list
/etc/dbus-1/system.d/dsme.conf
/etc/dsme
/etc/profile.d/reboot-via-dsme.sh
/usr/lib/startup/preinit/set_system_time
/usr/lib/systemd/system/dsme.service
/usr/lib/systemd/system/multi-user.target.wants/dsme.service
/usr/lib64/dsme
/usr/lib64/dsme/abootsettings.so
/usr/lib64/dsme/alarmtracker.so
/usr/lib64/dsme/batterytracker.so
/usr/lib64/dsme/bootreasonlogger.so
/usr/lib64/dsme/dbusautoconnector.so
/usr/lib64/dsme/dbusproxy.so
/usr/lib64/dsme/diskmonitor.so
/usr/lib64/dsme/emergencycalltracker.so
/usr/lib64/dsme/heartbeat.so
/usr/lib64/dsme/iphb.so
/usr/lib64/dsme/libabnormalexitwrapper.so
/usr/lib64/dsme/malf.so
/usr/lib64/dsme/processwd.so
/usr/lib64/dsme/pwrkeymonitor.so
/usr/lib64/dsme/runlevel.so
/usr/lib64/dsme/shutdownfeedback.so
/usr/lib64/dsme/startup.so
/usr/lib64/dsme/state.so
/usr/lib64/dsme/thermalflagger.so
/usr/lib64/dsme/thermalmanager.so
/usr/lib64/dsme/thermalsensor_generic.so
/usr/lib64/dsme/usbtracker.so
/usr/lib64/dsme/wlanloader.so
/usr/sbin/bootstate
/usr/sbin/dsme
/usr/sbin/dsme-server
/usr/sbin/dsmetemperature
/usr/sbin/dsmetool
/usr/sbin/getbootstate
/usr/sbin/waitfordsme
/usr/share/licenses/dsme-0.84.0
/usr/share/licenses/dsme-0.84.0/COPYING
/usr/share/licenses/dsme-0.84.0/copyright
/var/lib/dsme
/var/lib/dsme/alarm_queue_status
