# 4.5.0.21/RPMs.data/list/dummy_netd.list
/usr/lib/systemd/system/dummy_netd.service
/usr/lib/systemd/system/graphical.target.wants/dummy_netd.service
/usr/sbin/dummy_netd
