# 4.5.0.21/RPMs.data/list/jolla-sessions-qt5.list
/etc/pam.d/autologin
/etc/systemd/system/runlevel4.target
/sbin/preinit
/usr/lib/startup
/usr/lib/startup/autologind
/usr/lib/startup/create-home
/usr/lib/startup/get-initial-bootstate
/usr/lib/startup/hw-groups
/usr/lib/startup/init-done
/usr/lib/startup/limit-cpu
/usr/lib/startup/preinit
/usr/lib/startup/preinit/early.d
/usr/lib/startup/preinit/late.d
/usr/lib/startup/runlevel-change-done
/usr/lib/startup/start-autologin
/usr/lib/startup/validate-user
/usr/lib/startup/wait-for-file
/usr/lib/startup/wait-session-to-start
/usr/lib/systemd/system/actdead-pre.target
/usr/lib/systemd/system/actdead-pre.target.wants
/usr/lib/systemd/system/actdead-pre.target.wants/dbus.service
/usr/lib/systemd/system/actdead-pre.target.wants/dsme.service
/usr/lib/systemd/system/actdead-pre.target.wants/getty.target
/usr/lib/systemd/system/actdead-pre.target.wants/initial-bootstate.service
/usr/lib/systemd/system/actdead-pre.target.wants/mce.service
/usr/lib/systemd/system/actdead-pre.target.wants/oneshot-root.service
/usr/lib/systemd/system/actdead-pre.target.wants/sshd-keys.service
/usr/lib/systemd/system/actdead-pre.target.wants/sshd.socket
/usr/lib/systemd/system/actdead-pre.target.wants/systemd-ask-password-wall.path
/usr/lib/systemd/system/actdead-pre.target.wants/systemd-logind.service
/usr/lib/systemd/system/actdead-pre.target.wants/systemd-user-sessions.service
/usr/lib/systemd/system/actdead-pre.target.wants/usb-moded.service
/usr/lib/systemd/system/actdead.target
/usr/lib/systemd/system/actdead.target.wants
/usr/lib/systemd/system/actdead.target.wants/init-done.service
/usr/lib/systemd/system/actdead.target.wants/limit-cpu.service
/usr/lib/systemd/system/actdead.target.wants/runlevel-actdead-done.service
/usr/lib/systemd/system/actdead.target.wants/start-user-session.service
/usr/lib/systemd/system/autologin@.service
/usr/lib/systemd/system/create-home@.service
/usr/lib/systemd/system/display.service
/usr/lib/systemd/system/display.service.wants
/usr/lib/systemd/system/emergency-call.target
/usr/lib/systemd/system/emergency-call.target.wants
/usr/lib/systemd/system/emergency-call.target.wants/droid-hal-init.service
/usr/lib/systemd/system/emergency-call.target.wants/ofono.service
/usr/lib/systemd/system/emergency-call.target.wants/ohmd.service
/usr/lib/systemd/system/emergency-call.target.wants/pulseaudio.service
/usr/lib/systemd/system/graphical.target.wants/init-done.service
/usr/lib/systemd/system/graphical.target.wants/runlevel-user-done.service
/usr/lib/systemd/system/graphical.target.wants/start-user-session.service
/usr/lib/systemd/system/graphical.target.wants/wayland.path
/usr/lib/systemd/system/hw-groups@.service
/usr/lib/systemd/system/init-done.service
/usr/lib/systemd/system/initial-bootstate.service
/usr/lib/systemd/system/limit-cpu.service
/usr/lib/systemd/system/multi-user.target.wants/initial-bootstate.service
/usr/lib/systemd/system/quota@.service
/usr/lib/systemd/system/runlevel-actdead-done.service
/usr/lib/systemd/system/runlevel-user-done.service
/usr/lib/systemd/system/sdmount@.service
/usr/lib/systemd/system/start-user-session.service
/usr/lib/systemd/system/user@.service
/usr/lib/systemd/system/wayland.path
/usr/lib/systemd/system/wayland.service
/usr/lib/systemd/user/actdead-session-init-done.service
/usr/lib/systemd/user/booster-session.target.wants
/usr/lib/systemd/user/booster-session.target.wants/booster-silica-session.path
/usr/lib/systemd/user/booster-silica-session.path
/usr/lib/systemd/user/default.target
/usr/lib/systemd/user/post-actdead-session.target.wants
/usr/lib/systemd/user/post-actdead-session.target.wants/actdead-session-init-done.service
/usr/lib/systemd/user/post-user-session-init-done.service
/usr/lib/systemd/user/post-user-session.target.wants
/usr/lib/systemd/user/post-user-session.target.wants/user-session-init-done.service
/usr/lib/systemd/user/pre-user-session.target.wants
/usr/lib/systemd/user/pre-user-session.target.wants/dbus.service
/usr/lib/systemd/user/user-session-init-done.service
/usr/lib/tmpfiles.d/jolla-session-tmp.conf
/usr/lib/tmpfiles.d/sailfishos-cgroups.conf
/usr/libexec/sdmount
/usr/sbin/clear-device
/var/lib/environment
/var/lib/environment/nemo
/var/lib/environment/nemo/50-jolla-ui.conf
/var/lib/platform-updates
