# 4.5.0.21/RPMs.data/list/pulseaudio.list
/etc/dbus-1/system.d/pulseaudio-system.conf
/etc/pulse
/etc/pulse/client.conf
/etc/pulse/client.conf.d
/etc/pulse/client.conf.d/50-sfos.client.conf
/etc/pulse/daemon.conf
/etc/pulse/daemon.conf.d
/etc/pulse/daemon.conf.d/50-sfos.daemon.conf
/etc/pulse/default.pa
/etc/pulse/system.pa
/etc/security/limits.d/90-pulse.conf
/lib/udev/rules.d/90-pulseaudio.rules
/usr/bin/pa-info
/usr/bin/pacat
/usr/bin/pacmd
/usr/bin/pactl
/usr/bin/padsp
/usr/bin/pamon
/usr/bin/paplay
/usr/bin/parec
/usr/bin/parecord
/usr/bin/pasuspender
/usr/bin/pulseaudio
/usr/lib/systemd/system/pulseaudio.service
/usr/lib/systemd/user/pulseaudio.service
/usr/lib/systemd/user/pulseaudio.socket
/usr/lib/systemd/user/user-session.target.wants/pulseaudio.service
/usr/lib64/libpulse-mainloop-glib.so.0
/usr/lib64/libpulse-mainloop-glib.so.0.0.6
/usr/lib64/libpulse-simple.so.0
/usr/lib64/libpulse-simple.so.0.1.1
/usr/lib64/libpulse.so.0
/usr/lib64/libpulse.so.0.23.0
/usr/lib64/pulse-14.2
/usr/lib64/pulse-14.2/modules
/usr/lib64/pulse-14.2/modules/libalsa-util.so
/usr/lib64/pulse-14.2/modules/libbluez5-util.so
/usr/lib64/pulse-14.2/modules/libcli.so
/usr/lib64/pulse-14.2/modules/liboss-util.so
/usr/lib64/pulse-14.2/modules/libprotocol-cli.so
/usr/lib64/pulse-14.2/modules/libprotocol-http.so
/usr/lib64/pulse-14.2/modules/libprotocol-native.so
/usr/lib64/pulse-14.2/modules/libprotocol-simple.so
/usr/lib64/pulse-14.2/modules/librtp.so
/usr/lib64/pulse-14.2/modules/module-allow-passthrough.so
/usr/lib64/pulse-14.2/modules/module-alsa-card.so
/usr/lib64/pulse-14.2/modules/module-alsa-sink.so
/usr/lib64/pulse-14.2/modules/module-alsa-source.so
/usr/lib64/pulse-14.2/modules/module-always-sink.so
/usr/lib64/pulse-14.2/modules/module-always-source.so
/usr/lib64/pulse-14.2/modules/module-augment-properties.so
/usr/lib64/pulse-14.2/modules/module-bluetooth-discover.so
/usr/lib64/pulse-14.2/modules/module-bluetooth-policy.so
/usr/lib64/pulse-14.2/modules/module-bluez5-device.so
/usr/lib64/pulse-14.2/modules/module-bluez5-discover.so
/usr/lib64/pulse-14.2/modules/module-card-restore.so
/usr/lib64/pulse-14.2/modules/module-cli-protocol-tcp.so
/usr/lib64/pulse-14.2/modules/module-cli-protocol-unix.so
/usr/lib64/pulse-14.2/modules/module-cli.so
/usr/lib64/pulse-14.2/modules/module-combine-sink.so
/usr/lib64/pulse-14.2/modules/module-combine.so
/usr/lib64/pulse-14.2/modules/module-console-kit.so
/usr/lib64/pulse-14.2/modules/module-dbus-protocol.so
/usr/lib64/pulse-14.2/modules/module-default-device-restore.so
/usr/lib64/pulse-14.2/modules/module-detect.so
/usr/lib64/pulse-14.2/modules/module-device-manager.so
/usr/lib64/pulse-14.2/modules/module-device-restore.so
/usr/lib64/pulse-14.2/modules/module-echo-cancel.so
/usr/lib64/pulse-14.2/modules/module-filter-apply.so
/usr/lib64/pulse-14.2/modules/module-filter-heuristics.so
/usr/lib64/pulse-14.2/modules/module-gsettings.so
/usr/lib64/pulse-14.2/modules/module-hal-detect.so
/usr/lib64/pulse-14.2/modules/module-http-protocol-tcp.so
/usr/lib64/pulse-14.2/modules/module-http-protocol-unix.so
/usr/lib64/pulse-14.2/modules/module-intended-roles.so
/usr/lib64/pulse-14.2/modules/module-ladspa-sink.so
/usr/lib64/pulse-14.2/modules/module-loopback.so
/usr/lib64/pulse-14.2/modules/module-match.so
/usr/lib64/pulse-14.2/modules/module-mmkbd-evdev.so
/usr/lib64/pulse-14.2/modules/module-native-protocol-fd.so
/usr/lib64/pulse-14.2/modules/module-native-protocol-tcp.so
/usr/lib64/pulse-14.2/modules/module-native-protocol-unix.so
/usr/lib64/pulse-14.2/modules/module-null-sink.so
/usr/lib64/pulse-14.2/modules/module-null-source.so
/usr/lib64/pulse-14.2/modules/module-oss.so
/usr/lib64/pulse-14.2/modules/module-pipe-sink.so
/usr/lib64/pulse-14.2/modules/module-pipe-source.so
/usr/lib64/pulse-14.2/modules/module-position-event-sounds.so
/usr/lib64/pulse-14.2/modules/module-remap-sink.so
/usr/lib64/pulse-14.2/modules/module-remap-source.so
/usr/lib64/pulse-14.2/modules/module-rescue-streams.so
/usr/lib64/pulse-14.2/modules/module-role-cork.so
/usr/lib64/pulse-14.2/modules/module-role-ducking.so
/usr/lib64/pulse-14.2/modules/module-rtp-recv.so
/usr/lib64/pulse-14.2/modules/module-rtp-send.so
/usr/lib64/pulse-14.2/modules/module-rygel-media-server.so
/usr/lib64/pulse-14.2/modules/module-simple-protocol-tcp.so
/usr/lib64/pulse-14.2/modules/module-simple-protocol-unix.so
/usr/lib64/pulse-14.2/modules/module-sine-source.so
/usr/lib64/pulse-14.2/modules/module-sine.so
/usr/lib64/pulse-14.2/modules/module-stream-restore.so
/usr/lib64/pulse-14.2/modules/module-suspend-on-idle.so
/usr/lib64/pulse-14.2/modules/module-switch-on-connect.so
/usr/lib64/pulse-14.2/modules/module-switch-on-port-available.so
/usr/lib64/pulse-14.2/modules/module-systemd-login.so
/usr/lib64/pulse-14.2/modules/module-tunnel-sink-new.so
/usr/lib64/pulse-14.2/modules/module-tunnel-sink.so
/usr/lib64/pulse-14.2/modules/module-tunnel-source-new.so
/usr/lib64/pulse-14.2/modules/module-tunnel-source.so
/usr/lib64/pulse-14.2/modules/module-udev-detect.so
/usr/lib64/pulse-14.2/modules/module-virtual-sink.so
/usr/lib64/pulse-14.2/modules/module-virtual-source.so
/usr/lib64/pulse-14.2/modules/module-virtual-surround-sink.so
/usr/lib64/pulse-14.2/modules/module-volume-restore.so
/usr/lib64/pulseaudio
/usr/lib64/pulseaudio/libpulsecommon-14.2.so
/usr/lib64/pulseaudio/libpulsecore-14.2.so
/usr/lib64/pulseaudio/libpulsedsp.so
/usr/libexec/pulse
/usr/libexec/pulse/gsettings-helper
/usr/share/GConf/gsettings/pulseaudio.convert
/usr/share/bash-completion/completions/pacat
/usr/share/bash-completion/completions/pacmd
/usr/share/bash-completion/completions/pactl
/usr/share/bash-completion/completions/padsp
/usr/share/bash-completion/completions/paplay
/usr/share/bash-completion/completions/parec
/usr/share/bash-completion/completions/parecord
/usr/share/bash-completion/completions/pasuspender
/usr/share/bash-completion/completions/pulseaudio
/usr/share/glib-2.0/schemas/org.freedesktop.pulseaudio.gschema.xml
/usr/share/licenses/pulseaudio-14.2+git8
/usr/share/licenses/pulseaudio-14.2+git8/GPL
/usr/share/licenses/pulseaudio-14.2+git8/LGPL
/usr/share/licenses/pulseaudio-14.2+git8/LICENSE
/usr/share/locale/af/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/as/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/be/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/bg/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/bn_IN/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/ca/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/cs/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/da/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/de/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/de_CH/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/el/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/es/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/fi/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/fr/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/gl/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/gu/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/he/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/hi/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/hr/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/hu/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/id/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/it/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/ja/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/kk/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/kn/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/ko/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/lt/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/ml/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/mr/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/nl/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/nn/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/oc/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/or/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/pa/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/pl/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/pt/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/pt_BR/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/ru/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/sk/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/sr/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/sr@latin/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/sv/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/ta/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/te/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/tr/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/uk/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/zh_CN/LC_MESSAGES/pulseaudio.mo
/usr/share/locale/zh_TW/LC_MESSAGES/pulseaudio.mo
/usr/share/pulseaudio
/usr/share/pulseaudio/alsa-mixer
/usr/share/pulseaudio/alsa-mixer/paths
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-aux.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-dock-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-fm.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-front-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-headphone-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-headset-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-internal-mic-always.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-internal-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-linein.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-mic-line.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-mic.conf.common
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-rear-mic.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-tvtuner.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input-video.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-input.conf.common
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-headphones-2.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-headphones.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-lineout.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-mono.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-speaker-always.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-speaker.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output.conf.common
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-0.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-1.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-2.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-3.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-4.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-5.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-6.conf
/usr/share/pulseaudio/alsa-mixer/paths/hdmi-output-7.conf
/usr/share/pulseaudio/alsa-mixer/paths/iec958-stereo-input.conf
/usr/share/pulseaudio/alsa-mixer/paths/iec958-stereo-output.conf
/usr/share/pulseaudio/alsa-mixer/paths/steelseries-arctis-output-chat-common.conf
/usr/share/pulseaudio/alsa-mixer/paths/steelseries-arctis-output-game-common.conf
/usr/share/pulseaudio/alsa-mixer/paths/usb-gaming-headset-input.conf
/usr/share/pulseaudio/alsa-mixer/paths/usb-gaming-headset-output-mono.conf
/usr/share/pulseaudio/alsa-mixer/paths/usb-gaming-headset-output-stereo.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets
/usr/share/pulseaudio/alsa-mixer/profile-sets/audigy.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/cmedia-high-speed-true-hdaudio.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/default.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/dell-dock-tb16-usb-audio.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/force-speaker-and-int-mic.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/force-speaker.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/kinect-audio.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/maudio-fasttrack-pro.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-audio4dj.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-audio8dj.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-korecontroller.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-traktor-audio10.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-traktor-audio2.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-traktor-audio6.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/native-instruments-traktorkontrol-s4.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/sb-omni-surround-5.1.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/steelseries-arctis-common-usb-audio.conf
/usr/share/pulseaudio/alsa-mixer/profile-sets/usb-gaming-headset.conf
/usr/share/zsh/site-functions/_pulseaudio
