# 4.5.0.21/info/oneshot.list
/etc/oneshot.d
/etc/oneshot.d/group.d
/etc/oneshot.d/preinit
/etc/rpm/macros.oneshot
/etc/rpm/macros.oneshot-generic
/usr/bin/add-oneshot
/usr/bin/add-preinit-oneshot
/usr/bin/groupadd-user
/usr/bin/oneshot
/usr/bin/populate-oneshot-dir
/usr/lib/oneshot.d
/usr/lib/oneshot.d/groupadd-user.later
/usr/lib/systemd/system/graphical.target.wants/oneshot-root-late.service
/usr/lib/systemd/system/multi-user.target.wants/oneshot-root.service
/usr/lib/systemd/system/oneshot-root-late.service
/usr/lib/systemd/system/oneshot-root.service
/usr/lib/systemd/system/oneshot-user-late-privileged@.service
/usr/lib/systemd/system/oneshot-user-privileged@.service
/usr/lib/systemd/system/user@.service.d/oneshot-user-privileged.conf
/usr/lib/systemd/user/oneshot-user-late.service
/usr/lib/systemd/user/oneshot-user.service
/usr/lib/systemd/user/post-user-session.target.wants/oneshot-user-late.service
/usr/lib/systemd/user/pre-user-session.target.wants/oneshot-user.service
/usr/share/licenses/oneshot-0.6.12
/usr/share/licenses/oneshot-0.6.12/COPYING
/usr/share/user-managerd/remove.d/clean-oneshots.sh
