# 4.5.0.21/info/openssh-server.list
/etc/pam.d/sshd
/etc/ssh/sshd_config
/usr/lib/systemd/system/multi-user.target.wants/sshd-keys.service
/usr/lib/systemd/system/sshd-keys.service
/usr/lib/systemd/system/sshd.service
/usr/lib/systemd/system/sshd.socket
/usr/lib/systemd/system/sshd@.service
/usr/libexec/openssh/load_developer_profile
/usr/libexec/openssh/sftp-server
/usr/sbin/sshd
/usr/sbin/sshd-hostkeys
/var/empty/sshd
