# 4.5.0.21/info/sailfish-policy.list
/usr/lib/systemd/system/multi-user.target.wants/policies-setup.service
/usr/lib/systemd/system/policies-setup.service
/usr/lib64/libsailfishpolicy.so.0
/usr/lib64/libsailfishpolicy.so.0.4
/usr/lib64/libsailfishpolicy.so.0.4.17
/usr/lib64/qt5/qml/Sailfish/Policy
/usr/lib64/qt5/qml/Sailfish/Policy/libsailfishpolicyplugin.so
/usr/lib64/qt5/qml/Sailfish/Policy/plugins.qmltypes
/usr/lib64/qt5/qml/Sailfish/Policy/qmldir
/usr/libexec/policy-updater
/var/lib/policy
