# 4.5.0.21/info/sensorfw-qt5.list
/etc/dbus-1/system.d/sensorfw.conf
/etc/sensorfw
/usr/lib/systemd/system/graphical.target.wants/sensorfwd.service
/usr/lib/systemd/system/sensorfwd.service
/usr/lib64/libsensorclient-qt5.so.1
/usr/lib64/libsensorclient-qt5.so.1.0
/usr/lib64/libsensorclient-qt5.so.1.0.0
/usr/lib64/libsensordatatypes-qt5.so.1
/usr/lib64/libsensordatatypes-qt5.so.1.0
/usr/lib64/libsensordatatypes-qt5.so.1.0.0
/usr/lib64/libsensorfw-qt5.so.0
/usr/lib64/libsensorfw-qt5.so.0.9
/usr/lib64/libsensorfw-qt5.so.0.9.0
/usr/lib64/sensord-qt5
/usr/lib64/sensord-qt5/libaccelerometeradaptor-qt5.so
/usr/lib64/sensord-qt5/libaccelerometerchain-qt5.so
/usr/lib64/sensord-qt5/libaccelerometersensor-qt5.so
/usr/lib64/sensord-qt5/libalsadaptor-ascii-qt5.so
/usr/lib64/sensord-qt5/libalsadaptor-evdev-qt5.so
/usr/lib64/sensord-qt5/libalsadaptor-qt5.so
/usr/lib64/sensord-qt5/libalsadaptor-sysfs-qt5.so
/usr/lib64/sensord-qt5/libalssensor-qt5.so
/usr/lib64/sensord-qt5/libavgaccfilter-qt5.so
/usr/lib64/sensord-qt5/libcompasschain-qt5.so
/usr/lib64/sensord-qt5/libcompasssensor-qt5.so
/usr/lib64/sensord-qt5/libcoordinatealignfilter-qt5.so
/usr/lib64/sensord-qt5/libdeclinationfilter-qt5.so
/usr/lib64/sensord-qt5/libdownsamplefilter-qt5.so
/usr/lib64/sensord-qt5/libgyroadaptor-evdev-qt5.so
/usr/lib64/sensord-qt5/libgyroscopeadaptor-qt5.so
/usr/lib64/sensord-qt5/libgyroscopesensor-qt5.so
/usr/lib64/sensord-qt5/libhumidityadaptor-qt5.so
/usr/lib64/sensord-qt5/libhumiditysensor-qt5.so
/usr/lib64/sensord-qt5/libiiosensorsadaptor-qt5.so
/usr/lib64/sensord-qt5/libkbslideradaptor-qt5.so
/usr/lib64/sensord-qt5/liblidsensor-qt5.so
/usr/lib64/sensord-qt5/liblidsensoradaptor-evdev-qt5.so
/usr/lib64/sensord-qt5/libmagcalibrationchain-qt5.so
/usr/lib64/sensord-qt5/libmagcoordinatealignfilter-qt5.so
/usr/lib64/sensord-qt5/libmagnetometeradaptor-ascii-qt5.so
/usr/lib64/sensord-qt5/libmagnetometeradaptor-evdev-qt5.so
/usr/lib64/sensord-qt5/libmagnetometeradaptor-qt5.so
/usr/lib64/sensord-qt5/libmagnetometersensor-qt5.so
/usr/lib64/sensord-qt5/liborientationchain-qt5.so
/usr/lib64/sensord-qt5/liborientationinterpreter-qt5.so
/usr/lib64/sensord-qt5/liborientationsensor-qt5.so
/usr/lib64/sensord-qt5/libpressureadaptor-qt5.so
/usr/lib64/sensord-qt5/libpressuresensor-qt5.so
/usr/lib64/sensord-qt5/libproximityadaptor-ascii-qt5.so
/usr/lib64/sensord-qt5/libproximityadaptor-evdev-qt5.so
/usr/lib64/sensord-qt5/libproximityadaptor-qt5.so
/usr/lib64/sensord-qt5/libproximitysensor-qt5.so
/usr/lib64/sensord-qt5/librotationfilter-qt5.so
/usr/lib64/sensord-qt5/librotationsensor-qt5.so
/usr/lib64/sensord-qt5/libstepcountersensor-qt5.so
/usr/lib64/sensord-qt5/libtapadaptor-qt5.so
/usr/lib64/sensord-qt5/libtapsensor-qt5.so
/usr/lib64/sensord-qt5/libtemperatureadaptor-qt5.so
/usr/lib64/sensord-qt5/libtemperaturesensor-qt5.so
/usr/lib64/sensord-qt5/libtouchadaptor-qt5.so
/usr/sbin/sensorfwd
/usr/share/licenses/sensorfw-qt5-0.14.4
/usr/share/licenses/sensorfw-qt5-0.14.4/COPYING
