# 4.5.0.21/info/ssu.list
/etc/dbus-1/system.d/org.nemo.ssu.conf
/etc/ssu
/etc/zypp/credentials.d
/usr/bin/ssu
/usr/bin/ssud
/usr/lib/oneshot.d/ssu-update-repos
/usr/lib/systemd/system/dbus-org.nemo.ssu.service
/usr/lib64/libssu.so.1
/usr/lib64/libssu.so.1.0
/usr/lib64/libssu.so.1.0.24
/usr/lib64/ssu/libsandboxhook.so
/usr/libexec/zypp/plugins/urlresolver/ssu
/usr/share/dbus-1/system-services/org.nemo.ssu.service
/usr/share/licenses/ssu-1.0.24
/usr/share/licenses/ssu-1.0.24/COPYING
/usr/share/licenses/ssu-1.0.24/COPYING.BSD
/usr/share/licenses/ssu-1.0.24/COPYING.GPL
/usr/share/licenses/ssu-1.0.24/COPYING.LGPL
/usr/share/ssu
/usr/share/ssu/board-mappings.d
/usr/share/ssu/features.d
/usr/share/ssu/repos.d
