# 4.5.0.21/info/usb-moded.list
/etc/dbus-1/system.d/usb_moded.conf
/etc/modprobe.d/g_ether.conf
/etc/modprobe.d/usb_moded.conf
/etc/tmpfiles.d/usb-moded.conf
/etc/udhcpd.conf
/etc/usb-moded
/etc/usb-moded/10-usb-moded-defaults.ini
/etc/usb-moded/dyn-modes
/etc/usb-moded/run
/etc/usb-moded/usb-moded.ini
/usr/lib/systemd/system/basic.target.wants/usb-moded.service
/usr/lib/systemd/system/usb-moded.service
/usr/sbin/usb_moded
/usr/sbin/usb_moded_util
/usr/share/licenses/usb-moded-0.86.0+mer62
/usr/share/licenses/usb-moded-0.86.0+mer62/LICENSE
/usr/share/user-managerd/remove.d/usb_mode_user_clear.sh
/var/lib/usb-moded
/var/lib/usb-moded/usb-moded.ini
