# 4.5.0.21/RPMs.data/bluebinder.list
/usr/bin/droid/bluebinder_post.sh
/usr/bin/droid/bluebinder_wait.sh
/usr/lib/systemd/system/bluebinder.service
/usr/lib/systemd/system/graphical.target.wants/bluebinder.service
/usr/sbin/bluebinder
