# 4.5.0.21/RPMs.data/mce.list
/etc/dbus-1/system.d/mce.conf
/etc/mce
/etc/mce/10mce.ini
/etc/mce/20als-defaults.ini
/etc/mce/20debug-led.ini
/etc/mce/20evdev-types.ini
/etc/mce/20hybris-features.ini
/etc/mce/20hybris-led.ini
/etc/mce/20mce-radio-states.ini
/etc/tmpfiles.d/mce.conf
/usr/lib/systemd/system/mce.service
/usr/lib/systemd/system/multi-user.target.wants/mce.service
/usr/lib64/mce
/usr/lib64/mce/modules
/usr/lib64/mce/modules/alarm.so
/usr/lib64/mce/modules/audiorouting.so
/usr/lib64/mce/modules/battery-bme.so
/usr/lib64/mce/modules/battery-statefs.so
/usr/lib64/mce/modules/battery-udev.so
/usr/lib64/mce/modules/battery-upower.so
/usr/lib64/mce/modules/bluetooth.so
/usr/lib64/mce/modules/buttonbacklight.so
/usr/lib64/mce/modules/callstate.so
/usr/lib64/mce/modules/camera.so
/usr/lib64/mce/modules/charging.so
/usr/lib64/mce/modules/cpu-keepalive.so
/usr/lib64/mce/modules/display.so
/usr/lib64/mce/modules/doubletap.so
/usr/lib64/mce/modules/filter-brightness-als.so
/usr/lib64/mce/modules/fingerprint.so
/usr/lib64/mce/modules/inactivity.so
/usr/lib64/mce/modules/keypad.so
/usr/lib64/mce/modules/led.so
/usr/lib64/mce/modules/memnotify.so
/usr/lib64/mce/modules/mempressure.so
/usr/lib64/mce/modules/packagekit.so
/usr/lib64/mce/modules/powersavemode.so
/usr/lib64/mce/modules/proximity.so
/usr/lib64/mce/modules/radiostates.so
/usr/lib64/mce/modules/sensor-gestures.so
/usr/lib64/mce/modules/usbmode.so
/usr/sbin/dummy_compositor
/usr/sbin/mce
/usr/share/licenses/mce-1.114.1
/usr/share/licenses/mce-1.114.1/COPYING
/usr/share/licenses/mce-1.114.1/copyright
/var/lib/mce
