# 4.5.0.21/RPMs.data/nemo-qml-plugin-devicelock.list
/etc/dbus-1/system.d/org.nemomobile.devicelock.conf
/usr/lib/systemd/system/multi-user.target.wants/nemo-devicelock.socket
/usr/lib/systemd/system/nemo-devicelock.socket
/usr/lib64/libnemodevicelock.so.0
/usr/lib64/libnemodevicelock.so.0.3
/usr/lib64/libnemodevicelock.so.0.3.8
/usr/lib64/qt5/qml/org/nemomobile/devicelock
/usr/lib64/qt5/qml/org/nemomobile/devicelock/libnemodevicelockplugin.so
/usr/lib64/qt5/qml/org/nemomobile/devicelock/plugins.qmltypes
/usr/lib64/qt5/qml/org/nemomobile/devicelock/qmldir
/usr/share/licenses/nemo-qml-plugin-devicelock-0.3.8
/usr/share/licenses/nemo-qml-plugin-devicelock-0.3.8/LICENSE.LGPL
