# 4.5.0.21/RPMs.data/sailfish-device-encryption-service.list
/etc/crypttab
/etc/systemd/system/actdead.target.wants/jolla-actdead-charging.service
/etc/systemd/system/home-mount-settle.service.d/50-sailfish-home.conf
/etc/systemd/system/home.mount.d
/etc/systemd/system/home.mount.d/50-settle.conf
/etc/systemd/system/multi-user.target.d
/etc/systemd/system/multi-user.target.d/50-home.conf
/etc/systemd/system/systemd-user-sessions.service.d
/etc/systemd/system/systemd-user-sessions.service.d/50-home.conf
/usr/lib/systemd/system/aliendalvik.service.d/01-prevent-start.conf
/usr/lib/systemd/system/connman-vpn.service.d/01-prevent-start.conf
/usr/lib/systemd/system/connman.service.d/01-prevent-start.conf
/usr/lib/systemd/system/dbus-org.nemomobile.MmsEngine.service.d/01-prevent-start.conf
/usr/lib/systemd/system/dbus-org.nemomobile.provisioning.service.d/01-prevent-start.conf
/usr/lib/systemd/system/dbus-org.sailfishos.EncryptionService.service
/usr/lib/systemd/system/home-encryption-preparation.service
/usr/lib/systemd/system/home-mount-settle.service
/usr/lib/systemd/system/local-fs.target.wants/home-encryption-preparation.service
/usr/lib/systemd/system/mdm_proxy.service.d/01-prevent-start.conf
/usr/lib/systemd/system/mount-sd@.service.d/50-after-preparation.conf
/usr/lib/systemd/system/packagekit.service.d/01-home-mount.conf
/usr/lib/systemd/system/packagekit.service.d/01-prevent-start.conf
/usr/libexec/sailfish-encryption-service
/usr/share/dbus-1/system-services/org.sailfishos.EncryptionService.service
/usr/share/dbus-1/system.d/org.sailfishos.EncryptionService.conf
/usr/share/sailfish-device-encryption
/usr/share/sailfish-device-encryption/home-encryption-copy.sh
/usr/share/sailfish-device-encryption/home-encryption-finish.sh
/usr/share/sailfish-device-encryption/home-encryption-preparation.sh
/usr/share/sailfish-device-encryption/home-restore.sh
/var/lib/sailfish-device-encryption
